# conf/supervisord.conf
[supervisord]
nodaemon=true          # 核心：前台运行（容器主进程不能后台）
logfile=/var/log/supervisor/supervisord.log  # 主日志
pidfile=/var/run/supervisord.pid             # 进程ID文件
childlogdir=/var/log/supervisor              # 子程序日志目录

# 配置第一个程序：（必须前台运行）
[program:spider]
command=/usr/local/bin/fly-go --service=server  
autostart=true                            # 随Supervisord启动
autorestart=true                          # 进程崩溃自动重启
stdout_logfile=/var/log/supervisor/nginx.log  # 标准输出日志
stderr_logfile=/var/log/supervisor/nginx.err.log  # 错误日志
user=root                                 # 运行用户

# 配置第二个程序：（必须前台运行）
[program:server]
command=/usr/local/bin/fly-go  --service=spider 
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/sshd.log
stderr_logfile=/var/log/supervisor/sshd.err.log
user=root

# 如需添加更多程序，复制[program:xxx]块并修改即可（如php-fpm、自定义脚本）