# build/fly-go.service.template
[Unit]
Description=fly-go service
After=network.target

[Service]
# Replace the ExecStart path below with the absolute path to your binary.
ExecStart=/usr/local/bin/fly-go
Restart=on-failure
RestartSec=5s
User=root
Environment=GIN_MODE=release

[Install]
WantedBy=multi-user.target
